# xdg-user-dirs/media/files/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin mc --noclear %I 38400 linux
