# xdg-user-dirs/media/files/etc/rc.local
#!/bin/bash

echo 10 > /proc/sys/kernel/panic
cmd=$(lsusb | grep 'SCM Microsystems, Inc.' | awk -F ':' '{ print $1 }' | awk '{ printf("chown root:media /dev/bus/usb/%s/%s", $2, $4) }')
$cmd

cpupower frequency-set --governor performance --min 1Ghz --max 2Ghz

exit 0
