# xdg-user-dirs/media/files/etc/rc.local
#!/bin/bash

echo 10 > /proc/sys/kernel/panic
cpupower frequency-set --max 1.8Ghz

sudo -u mc /home/mc/xdg-user-dirs/media/bin/startup.sh &

exit 0
