# xdg-user-dirs/media/files/etc/init/block-device-added-ssd-array.conf
start on block-device-added DEVNAME=/dev/md1p1

task
emits mounted-ssd-array

script

    date +"%Y/%m/%d %H:%M:%S.%N START" >> /tmp/$UPSTART_JOB

    mount -o relatime,nouser_xattr /dev/md1p1 /mnt/ssd_array
    initctl emit mounted-ssd-array

    date +"%Y/%m/%d %H:%M:%S.%N END" >> /tmp/$UPSTART_JOB

end script
